# cfg/systemd/cmschecker.service
[Unit]
Description=CMS Heartbeat Monitor Service
After=cmslog.service

[Service]
Type=simple
User=demon
WorkingDirectory=/home/demon/cms_venv/bin
ExecStart=/home/demon/cms_venv/bin/cmsChecker
Restart=on-failure

[Install]
WantedBy=multi-user.target
