# cfg/systemd/cmsscoring.service
[Unit]
Description=CMS Scoring Service
After=cmslog.service

[Service]
Type=simple
User=demon
WorkingDirectory=/home/demon/cms_venv/bin
ExecStart=/home/demon/cms_venv/bin/cmsScoringService
Restart=on-failure

[Install]
WantedBy=multi-user.target
